# config/mystic-monitor.service
[Unit]
Description=Mystic Monitor Background ML Daemon
After=network.target

[Service]
Type=simple
User=root
Group=mystic
WorkingDirectory=/opt/mystic_monitor
ExecStart=/usr/bin/python3 /opt/mystic_monitor/mystic_daemon.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Systemd Integration & Run Directory
RuntimeDirectory=mystic
RuntimeDirectoryMode=0750

# Shutdown Behavior
KillSignal=SIGTERM
TimeoutStopSec=10

# Hardening & Sandboxing Knobs
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
ReadWritePaths=/var/log/mystic-anomalies.log /run/mystic

[Install]
WantedBy=multi-user.target
